# supervisord-studio.conf
[program:studio]
command=/usr/bin/python3 /home/unsloth/.local/lib/python3.12/site-packages/studio/backend/run.py --host 0.0.0.0 --port 8000
user=unsloth
directory=/workspace
autostart=true
autorestart=true
environment=HOME="/home/unsloth",USER="unsloth"
stdout_logfile=/var/log/studio/access.log
stderr_logfile=/var/log/studio/error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
