# supervisord-studio.conf
[program:studio]
command=/home/unsloth/.local/bin/unsloth studio -H 0.0.0.0 -p 8000 -q
user=unsloth
directory=/workspace
autostart=true
autorestart=true
environment=HOME="/home/unsloth",USER="unsloth"
stdout_logfile=/var/log/studio/access.log
stderr_logfile=/var/log/studio/error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
