# supervisord-studio.conf
[program:studio]
command=bash -c "source /home/unsloth/.unsloth/studio/unsloth_studio/bin/activate && PYTHONPATH=/usr/local/lib/python3.12/dist-packages unsloth studio -H 0.0.0.0 -p 8000 -q"
user=unsloth
directory=/workspace
autostart=true
autorestart=true
environment=HOME="/home/unsloth",USER="unsloth",PATH="/home/unsloth/.unsloth/studio/unsloth_studio/bin:/home/unsloth/.local/bin:/usr/local/cuda/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
stdout_logfile=/var/log/studio/access.log
stderr_logfile=/var/log/studio/error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
